# handle_dpms.service
Description=DPMS event handler

[Service]
ExecStart=/usr/local/libexec/handle_dpms

[Install]
WantedBy=multi-user.target
