# handle_dpms.service
[Unit]
Description=DPMS event handler

[Service]
ExecStart=/usr/libexec/handle_dpms

[Install]
WantedBy=multi-user.target
